# docker/rootfs/etc/s6-overlay/s6-rc.d/init/run
#!/command/with-contenv sh

# Copy/update td shim
cat "$TYPEDAEMON_MODULE/skel/bin/td" > "$TYPEDAEMON_CONFIG/bin/td"
