# docker/rootfs/etc/s6-overlay/s6-rc.d/sshd_socket/run
#!/command/execlineb -P
socat -lf socat.log UNIX-LISTEN:$TYPEDAEMON_CONFIG/.typedaemon/ssh_sock,reuseaddr,fork TCP:localhost:22