# docker/rootfs/etc/s6-overlay/s6-rc.d/sshd/run
#!/command/execlineb -P
sshd < /dev/null
