# docker/rootfs/etc/s6-overlay/s6-rc.d/init-bin/run
#!/command/with-contenv sh

mkdir -p "$TYPEDAEMON_CONFIG/.typedaemon/bin"

cp $(which ssh) "$TYPEDAEMON_CONFIG/.typedaemon/bin/ssh"
cp $(which socat) "$TYPEDAEMON_CONFIG/.typedaemon/bin/socat"
